# TestDocker/VNC/supervisord.conf
[supervisord]
nodaemon=true

[program:vnc]
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24
user=vncuser
autostart=true
autorestart=true